# src/scripts/com.user.macosgithubprojects.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.macosgithubprojects</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/clm/Documents/GitHub/PROJECTS/Macos_GithubProjects/.venv/bin/python3</string>
        <string>/Users/clm/Documents/GitHub/PROJECTS/Macos_GithubProjects/src/macos_githubprojects/menu_app_improved.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/macosgithubprojects.stdout</string>
    <key>StandardErrorPath</key>
    <string>/tmp/macosgithubprojects.stderr</string>
</dict>
</plist>